# themes/luci-theme-firefly/root/etc/uci-defaults/luci-theme-firefly
#!/bin/sh
uci batch <<-EOF
        set luci.themes.FireFly=/luci-static/firefly
	set luci.main.mediaurlbase=/luci-static/firefly
        commit luci
EOF
exit 0
